# ot-reddish.service
# Reference unit for “pi@/home/pi/OT-Reddish” — install_usb.sh rewrites this for the real user and paths.
[Unit]
Description=OT-Reddish Scanner (Kiosk)
After=network.target graphical.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/OT-Reddish
Environment=HOME=/home/pi
Environment=DISPLAY=:0
Environment=OT_REDDISH_PI_MODE=1
Environment=OT_REDDISH_KIOSK=1
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/home/pi/OT-Reddish/venv/bin/python /home/pi/OT-Reddish/app.py
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
